# systemd/email-learner-worker.service
[Unit]
Description=Agent Mail 24/7 - Action Worker (consomme action_queue vers Gmail API)
Documentation=https://www.explodev.fr
After=postgresql.service network-online.target
Wants=network-online.target

[Service]
Type=simple
User=eddie
Group=eddie
WorkingDirectory=/home/eddie/email-learner
Environment="PATH=/home/eddie/email-learner/venv/bin:/usr/local/bin:/usr/bin"
EnvironmentFile=/home/eddie/email-learner/configs/.env
ExecStart=/home/eddie/email-learner/venv/bin/python -m src.action_worker
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=email-learner-worker

# Securite : pas d'acces root, pas de privilege escalation
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/eddie/email-learner/configs

# Resource limits
LimitNOFILE=65536
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
